# generic/rootfs/etc/services.d/home-assistant/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Home Assistant service
# ==============================================================================
cd /config
exec python3 -m homeassistant --config /config
