# landingpage/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Nginx service
# ==============================================================================
exec nginx
